# crontab
@daily python -u /app/whois.py >> /logfile
@reboot python -u /app/whois.py >> /logfile
